# linux-files/rc.local
#!/bin/sh -e
# This file was created by SysprePi.
# More information at: https://github.com/bobeirasa/SysprePi

/home/pi/SysprePi/pi-bootstrapper/bootstrap.sh >> /tmp/rclocal

exit 0
